# launchd/com.contractwatch.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  ContractWatch monthly bulk-refresh launchd job (macOS).

  Install (one-time):
    1. Copy this file to ~/Library/LaunchAgents/com.contractwatch.plist
    2. Edit the two WorkingDirectory / ProgramArguments paths to your absolute
       checkout path (must be absolute; launchd does not expand ~).
    3. launchctl load ~/Library/LaunchAgents/com.contractwatch.plist

  monthly_scan.sh runs the full bulk-archive refresh (download all FY15-FY26
  archives, rebuild the DB, reflag, export JSON), then optionally pushes to
  Cloudflare Pages. Expect ~26 minutes wall time and ~25 GB of free disk.

  StartCalendarInterval below fires the 8th of every month at 07:00 local.
  USASpending publishes the monthly bulk archive snapshot on the 5th-6th,
  so the 8th gives a small buffer for upstream data to settle.
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.contractwatch</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/REPLACE/WITH/ABSOLUTE/PATH/TO/contractwatch/monthly_scan.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/REPLACE/WITH/ABSOLUTE/PATH/TO/contractwatch</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Day</key>
    <integer>8</integer>
    <key>Hour</key>
    <integer>7</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>RunAtLoad</key>
  <false/>

  <key>StandardOutPath</key>
  <string>/tmp/contractwatch.out.log</string>

  <key>StandardErrorPath</key>
  <string>/tmp/contractwatch.err.log</string>

  <key>ProcessType</key>
  <string>Background</string>
</dict>
</plist>
